# router/overlay/upper/usr/lib/opkg/info/stunnel.list
/etc/init.d/stunnel
/usr/bin/stunnel
/usr/lib/stunnel/libstunnel.so
/etc/config/stunnel
/etc/stunnel/stunnel.conf
