# router/etc/hotplug.d/iface/01-disable-modem-retry
[ -n "$(echo $INTERFACE|grep modem)" ] && [ "ifup-failed" = "$ACTION" ] && {
	. /lib/netifd/netifd-proto.sh
	#disable modem auto retry,waiting gl-monitor try again
	proto_block_restart "$INTERFACE"
}
